# docker/supervisord.conf
[supervisord]

[program:rackup]
command=rackup

[program:foreman]
command=bundle exec foreman start
